# support/moonscale.service
[Unit]
Description=Moonscale service

[Service]
User=pi
WorkingDirectory=/home/pi
ExecStart=/home/pi/moonscale/start.sh

[Install]
WantedBy=multi-user.target
